# weekly.cron
0 2 * * 1 cd /home/csg0026/project/datagetta && ./cron-script.sh
